# src/program-influxdb.conf
[program:influxdb]
command = /usr/local/bin/launch-influxdb.sh

redirect_stderr = true

stdout_logfile = /var/log/influxdb/influxdb.log
stdout_logfile_maxbytes = 5MB
